# services/hedge-performance.service
[Unit]
Description=HEDGE Performance Reporter
After=network.target

[Service]
Type=oneshot
User=teckz
Group=teckz
WorkingDirectory=/home/teckz/Documents/HEDGE
Environment="PATH=/home/teckz/Documents/HEDGE/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/home/teckz/Documents/HEDGE/.env
ExecStart=/home/teckz/Documents/HEDGE/venv/bin/python -m core.services.performance_reporter
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

