# services/hedge-performance.timer
[Unit]
Description=HEDGE Performance Reporter (daily)
Requires=hedge-performance.service

[Timer]
OnCalendar=Mon..Fri 09:30
Unit=hedge-performance.service
Persistent=true

[Install]
WantedBy=timers.target

