# services/hedge-health.timer
[Unit]
Description=HEDGE Health Monitor (nightly)
Requires=hedge-health.service

[Timer]
OnCalendar=*-*-* 23:55
Unit=hedge-health.service
Persistent=true

[Install]
WantedBy=timers.target

